# scripts/com.user.cpukiller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.cpukiller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/adamkruger/Personal/consensus-app/scripts/cpu-killer.sh</string>
        <string>250</string>
        <string>5</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/adamkruger/cpu-killer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/adamkruger/cpu-killer-error.log</string>
</dict>
</plist>

